# private_dot_config/systemd/user/i3doro.service
[Unit]
Description=Pomodoro daemon

[Service]
Type=simple
ExecStart=/home/alan/git/i3doro/i3doro.py

[Install]
WantedBy=default.target
