# private_dot_config/systemd/user/picom.service
[Unit]
Description=Picom compositor
[Service]
ExecStart=/usr/bin/picom
Restart=always
RestartSec=3
[Install]
WantedBy=default.target
